# tpadrestart.service
[Unit]
Description=Touch pad Restart

[Service]
ExecStart=/usr/local/bin/touchpadrestart

[Install]
WantedBy=multi-user.target
